# infrastructure/ansible/templates/openclaw-node-exporter.service.j2
[Unit]
Description=Prometheus node_exporter for OpenClaw host metrics
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/node_exporter --web.listen-address=:9100
Restart=on-failure

[Install]
WantedBy=multi-user.target
