# infrastructure/ansible/templates/k8docker/node-exporter.service.j2
[Unit]
Description=Prometheus node_exporter for k8docker host metrics
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/node_exporter \
  --web.listen-address=:9100 \
  --collector.textfile.directory={{ k8docker_db_metrics_dir }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
